# deploy/systemd/cmc-author-dashboard-web.service
[Unit]
Description=CMC Author Dashboard Web
After=network.target

[Service]
Type=simple
User=thana
WorkingDirectory=/opt/cmc-author-dashboard
Environment=NODE_ENV=production
Environment=HOME=/home/thana
Environment=PATH=/home/thana/.npm-global/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=/opt/cmc-author-dashboard/.env
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
